# image-builder/sovereign/layer/sovereign-proof.rootfs-overlay/etc/systemd/system/sovereign-pihole.service
[Unit]
Description=Configure and start the Sovereign Pi-hole appliance
Requires=sovereign-pihole-import.service docker.service
After=sovereign-pihole-import.service docker.service network-online.target
Wants=network-online.target
RequiresMountsFor=/data
ConditionPathIsMountPoint=/data

[Service]
Type=oneshot
ExecStart=/usr/lib/sovereign/start-pihole
RemainAfterExit=yes
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
