# image-builder/sovereign/layer/sovereign-proof.rootfs-overlay/etc/systemd/system/sovereign-credentials-console.service
[Unit]
Description=Show Sovereign first-boot credentials on the physical console
Requires=sovereign-local-access.service
After=sovereign-local-access.service getty@tty1.service
ConditionPathExists=/data/sovereign/secrets/pihole-admin-password

[Service]
Type=oneshot
ExecStart=/usr/lib/sovereign/show-pihole-credentials
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1
TTYReset=no
TTYVHangup=no
TTYVTDisallocate=no
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
