# image-builder/sovereign/layer/sovereign-proof.rootfs-overlay/etc/systemd/system/sovereign-pihole-artifact.service
[Unit]
Description=Verify the embedded official Pi-hole ARM64 artifact
RequiresMountsFor=/data
After=local-fs.target sovereign-proof.service
ConditionPathIsMountPoint=/data

[Service]
Type=oneshot
ExecStart=/usr/lib/sovereign/verify-pihole-artifact
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
