# image-builder/sovereign/layer/sovereign-proof.rootfs-overlay/etc/systemd/system/sovereign-local-access.service
[Unit]
Description=Verify Sovereign local HTTP access
Requires=sovereign-pihole.service nginx.service
After=sovereign-pihole.service nginx.service
Wants=avahi-daemon.service
RequiresMountsFor=/data
ConditionPathIsMountPoint=/data
StartLimitIntervalSec=300
StartLimitBurst=6

[Service]
Type=oneshot
ExecStart=/usr/lib/sovereign/verify-local-access
RemainAfterExit=yes
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
