# image-builder/sovereign/layer/sovereign-proof.rootfs-overlay/etc/systemd/system/sovereign-proof.service
[Unit]
Description=Initialize the Sovereign proof data directory
RequiresMountsFor=/data
After=local-fs.target
ConditionPathIsMountPoint=/data

[Service]
Type=oneshot
ExecStart=/usr/lib/sovereign/proof-init
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

