# image-builder/sovereign/layer/sovereign-proof.rootfs-overlay/etc/systemd/system/containerd.service.d/10-sovereign-data.conf
[Unit]
RequiresMountsFor=/data
After=local-fs.target
ConditionPathIsMountPoint=/data

[Service]
ExecStartPre=/usr/bin/install -d -m 0711 /data/containerd
