# image-builder/sovereign/layer/sovereign-proof.rootfs-overlay/etc/systemd/system/nginx.service.d/10-sovereign-pihole.conf
[Unit]
Wants=sovereign-pihole.service sovereign-console.service
After=sovereign-pihole.service sovereign-console.service
