# image-builder/sovereign/layer/sovereign-proof.rootfs-overlay/etc/systemd/system/nginx.service.d/10-sovereign-pihole.conf
[Unit]
Requires=sovereign-pihole.service
After=sovereign-pihole.service
