# image-builder/sovereign/layer/sovereign-proof.rootfs-overlay/etc/systemd/system/iwd.service.d/10-sovereign-provisioning.conf
[Unit]
After=sovereign-imager-provision.service
Wants=sovereign-imager-provision.service
